# infra/upend@.service
# systemd template unit — one instance per service
# Usage: systemctl enable upend@api upend@claude
[Unit]
Description=upend %i service
After=network.target

[Service]
Type=simple
User=ec2-user
WorkingDirectory=/opt/upend
EnvironmentFile=/opt/upend/.env
ExecStart=/usr/local/bin/bun services/%i/index.ts
Restart=always
RestartSec=2
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
